# agate-dist/src/main/deb/debian/postinst
#!/bin/sh
# postinst script for agate
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)

    	RET=""
      db_get agate/admin_password

      if [ -z "$2" ] && [ "$RET" != "" ] ; then
        echo Updating administrator password.
        # hash the password and escape the possible '$/&' characters
        adminpw=$(echo -n $RET | xargs java -jar /usr/share/agate-*/tools/lib/obiba-password-hasher-*-cli.jar | sed -e 's/\([\$\/\/]\)/\\\1/g' | cut -d ' ' -f 1)
        tfile=`mktemp`
        if [ ! -f "$tfile" ]; then
          return 1
        fi
        cat /etc/agate/shiro.ini | sed -e "s,^administrator\s*=.*\,,administrator=$adminpw\,," > $tfile && \
          mv $tfile /etc/agate/shiro.ini
        # clear password in db_get so it won't be re-encoded during next upgrade
        db_set agate/admin_password ""
      fi

    	# Create agate user if it doesn't exist.
      if ! id agate > /dev/null 2>&1 ; then
        adduser --system --home /var/lib/agate --no-create-home --disabled-password agate
      fi

      # Agate file structure on Debian
      # /etc/agate: configuration
      # /usr/share/agate: executable
      # /var/lib/agate: data runtime
      # /var/log: logs

      rm -f /usr/share/agate
      ln -s /usr/share/agate-* /usr/share/agate

      if [ ! -e /var/lib/agate/conf ] ; then
        ln -s /etc/agate /var/lib/agate/conf
      fi

      chown -R agate:adm /var/lib/agate /var/log/agate /etc/agate
      chmod -R 750      /var/lib/agate /var/log/agate /etc/agate
      chmod +x /usr/share/agate/tools/shiro-hasher
      find /etc/agate/ -type f -print0 | xargs -0 chmod 640

      if [ -e /etc/init.d/agate ]; then
        update-rc.d agate defaults >/dev/null
        rm /etc/init.d/agate
      fi
      systemctl daemon-reload
      systemctl enable agate
      systemctl start agate
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
